# openwrt/files/etc/uci-defaults/99-nikki
uci batch << EOF
set nikki.config.start_delay='3'
set nikki.config.mixin='0'
set nikki.config.fast_reload='0'
set nikki.mixin.dns_mode='redir-host'
set nikki.mixin.fake_ip_range='28.0.0.1/8'
set nikki.mixin.ui_url='https://raw.ihtw.moe/https://github.com/Zephyruso/zashboard/archive/refs/heads/gh-pages.zip'
commit
EOF